# daemon/com.sogos.claude-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sogos.claude-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__DAEMON_PYTHON__</string>
        <string>__DAEMON_DIR__/daemon.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__DAEMON_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/claude-daemon/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/claude-daemon/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
